# configs/keepalived.service
# routedns-ingress — systemd unit for Keepalived (source installs)
# Distro packages ship their own unit; this is used when building Keepalived from source.

[Unit]
Description=Keepalive Daemon (LVS and VRRP)
After=network-online.target
Wants=network-online.target

[Service]
Type=forking
PIDFile=/run/keepalived.pid
KillMode=process
EnvironmentFile=-/etc/default/keepalived
# This build has no systemd-notify support, so Type=forking is used.
# Do NOT pass --dont-fork here: it keeps keepalived in the foreground and
# systemd would wait for the (never-exiting) parent until the start timeout.
ExecStart=/usr/sbin/keepalived $DAEMON_ARGS
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
